# backup/backup.service
[Unit]
Description="run toolbox-backup"

[Service]
ExecStart=/usr/bin/toolbox-backup
